# backend.service
# Systemd service файл для backend (опционально, если не используете docker-compose)
# Разместите в /etc/systemd/system/backend.service
# ВАЖНО: Замените ~/plannerAI на полный путь (например /home/user/plannerAI)
# Использование:
#   sudo systemctl daemon-reload
#   sudo systemctl enable backend
#   sudo systemctl start backend

[Unit]
Description=AI Planner Backend Service
After=network.target

[Service]
Type=simple
User=www-data
# Замените ~/plannerAI на полный путь, например /home/user/plannerAI
WorkingDirectory=/home/user/plannerAI
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/user/plannerAI/.env
ExecStart=/usr/local/bin/docker compose up
ExecStop=/usr/local/bin/docker compose down
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

